# Terminal/etc/profile.d/goenv

export GOROOT=/usr/local/go
export GOPATH=$HOME/go
export GOBIN=/usr/local/go/bin

export PATH=$PATH:$GOROOT/bin:${PATH}:${GOPATH}/bin
